# vagrant/spider_host/artifacts/clf-spider-host.conf
description "CLF Spider Host"
author "dave@cloudfeaster.com"

start on runlevel [2345]
stop on runlevel [!2345]

setuid clf
setgid clf

env DISPLAY=:99

respawn

exec spider_host --log=info
